# contrib/systemd/cattenbak.timer
# systemd timer for cattenbak
[Unit]
Description=Run cattenbak hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
